# contrib/systemd/alfis.service
[Unit]
Description=alfis
Wants=network.target
Wants=alfis-default-config.service
After=network.target
After=alfis-default-config.service

[Service]
User=alfis
Group=alfis

# Memory limits to prevent OOM kills and memory leaks
# MemoryHigh: soft limit - systemd will throttle the process if exceeded
# MemoryMax: hard limit - process will be killed if exceeded
# Default values: 2GB (can be adjusted based on system resources)
MemoryHigh=2G
MemoryMax=2G

ProtectHome=true
ProtectSystem=true

SecureBits=keep-caps
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

SyslogIdentifier=alfis
WorkingDirectory=/var/lib/alfis
ExecStart=/usr/bin/alfis -n -c /etc/alfis.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
